# bin/set-theme@.service
[Unit]
Description=Set theme to %i

[Service]
Type=oneshot
ExecStart=/home/max/dotfiles/bin/set-theme.sh %i
